# web/etc/cron.hourly/lyx2rfc
#!/bin/sh

find /tmp -name "php*" -amin +10 -exec rm {} \;
